# proyecto/systemd/daemon/lab.target
[Unit]
Description=Proyecto bases 2 daemon y carga de dependencias
Requires= mongod.service lab.service flask.service
